# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
loglevel=info
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=5
startretries=10

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10

[program:cron]
command=crond -f -d 8
autostart=true
autorestart=true
